# jobs/nfsclient/templates/data/crontab
* * * * * root /var/vcap/jobs/nfsclient/bin/ctl mount